# cloud/linode/ba.service
[Unit]
Description=baba web server

[Service]
ExecStart=/a/ba.sh

[Install]
WantedBy=multi-user.target

